# clubalpin_rss.service
[Unit]
Description=Club Alpin RSS Feed Server
After=network.target

[Service]
Type=simple
User=pbeaudequin
WorkingDirectory=/Users/pbeaudequin/clubalpin_rss
Environment=PYTHONUNBUFFERED=1
ExecStart=/Users/pbeaudequin/clubalpin_rss/venv/bin/python /Users/pbeaudequin/clubalpin_rss/scrape_to_rss.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
